# etc/supervisord.d/docker.ini
[program:docker]
user=root
autostart=true
autorestart=true
command=/usr/sbin/wrapdocker
